# turris-news.cron
MAILTO=""
0	*	*	*	*	root	/usr/bin/turris-news cron
